# desktop_arch_hyprland/etc/systemd/system/getty@tty1.service.d/autologin.conf
[Service]
ExecStart=
ExecStart=-/sbin/agetty -o '-p -f -- \\u' --noclear --autologin YOU_NEED_TO_EDIT_THIS_USERNAME %I $TERM